# s6/rclone/finish
#!/bin/bash -x

/bin/fusermount -u /mnt/rclone
