# s6/.s6-svscan/finish
#!/bin/bash -x

/bin/umount -f /mnt/unionfs

/bin/umount -f /mnt/rclone
